# systemd/wifi-setup.service
[Unit]
Description=BLE WiFi Setup Service
After=bluetooth.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/path/to/your/wifi-setup

[Install]
WantedBy=multi-user.target
